# autostart/quickshell.desktop
[Desktop Entry]
Type=Application
Name=Quickshell
Exec=quickshell -c caelestia
